# profile.d/cygcac.sh
export OPENSSL_CONF=/usr/local/ssl/openssl.cnf
